# config/supervisor.conf
; Notes:
; priority=1 --> Lower priorities indicate programs that start first and shut down last
; killasgroup=true --> send kill signal to child processes too

[program:ecommerce-frappe-frappe-web]
command=/home/yabesh/Downloads/Office/frappe/ecommerce-frappe/env/bin/gunicorn -b 127.0.0.1:8002 -w 17 --max-requests 5000 --max-requests-jitter 500 -t 120 frappe.app:application --preload
priority=4
autostart=true
autorestart=true
stdout_logfile=/home/yabesh/Downloads/Office/frappe/ecommerce-frappe/logs/web.log
stderr_logfile=/home/yabesh/Downloads/Office/frappe/ecommerce-frappe/logs/web.error.log
user=yabesh
directory=/home/yabesh/Downloads/Office/frappe/ecommerce-frappe/sites


[program:ecommerce-frappe-frappe-schedule]
command=/home/yabesh/.local/bin/bench schedule
priority=3
autostart=true
autorestart=true
stdout_logfile=/home/yabesh/Downloads/Office/frappe/ecommerce-frappe/logs/schedule.log
stderr_logfile=/home/yabesh/Downloads/Office/frappe/ecommerce-frappe/logs/schedule.error.log
user=yabesh
directory=/home/yabesh/Downloads/Office/frappe/ecommerce-frappe

[program:ecommerce-frappe-frappe-default-worker]
command=/home/yabesh/.local/bin/bench worker --queue default
priority=4
autostart=true
autorestart=true
stdout_logfile=/home/yabesh/Downloads/Office/frappe/ecommerce-frappe/logs/worker.log
stderr_logfile=/home/yabesh/Downloads/Office/frappe/ecommerce-frappe/logs/worker.error.log
user=yabesh
stopwaitsecs=1560
directory=/home/yabesh/Downloads/Office/frappe/ecommerce-frappe
killasgroup=true
numprocs=1
process_name=%(program_name)s-%(process_num)d

[program:ecommerce-frappe-frappe-short-worker]
command=/home/yabesh/.local/bin/bench worker --queue short
priority=4
autostart=true
autorestart=true
stdout_logfile=/home/yabesh/Downloads/Office/frappe/ecommerce-frappe/logs/worker.log
stderr_logfile=/home/yabesh/Downloads/Office/frappe/ecommerce-frappe/logs/worker.error.log
user=yabesh
stopwaitsecs=360
directory=/home/yabesh/Downloads/Office/frappe/ecommerce-frappe
killasgroup=true
numprocs=1
process_name=%(program_name)s-%(process_num)d

[program:ecommerce-frappe-frappe-long-worker]
command=/home/yabesh/.local/bin/bench worker --queue long
priority=4
autostart=true
autorestart=true
stdout_logfile=/home/yabesh/Downloads/Office/frappe/ecommerce-frappe/logs/worker.log
stderr_logfile=/home/yabesh/Downloads/Office/frappe/ecommerce-frappe/logs/worker.error.log
user=yabesh
stopwaitsecs=1560
directory=/home/yabesh/Downloads/Office/frappe/ecommerce-frappe
killasgroup=true
numprocs=1
process_name=%(program_name)s-%(process_num)d






[program:ecommerce-frappe-redis-cache]
command=/usr/bin/redis-server /home/yabesh/Downloads/Office/frappe/ecommerce-frappe/config/redis_cache.conf
priority=1
autostart=true
autorestart=true
stdout_logfile=/home/yabesh/Downloads/Office/frappe/ecommerce-frappe/logs/redis-cache.log
stderr_logfile=/home/yabesh/Downloads/Office/frappe/ecommerce-frappe/logs/redis-cache.error.log
user=yabesh
directory=/home/yabesh/Downloads/Office/frappe/ecommerce-frappe/sites

[program:ecommerce-frappe-redis-queue]
command=/usr/bin/redis-server /home/yabesh/Downloads/Office/frappe/ecommerce-frappe/config/redis_queue.conf
priority=1
autostart=true
autorestart=true
stdout_logfile=/home/yabesh/Downloads/Office/frappe/ecommerce-frappe/logs/redis-queue.log
stderr_logfile=/home/yabesh/Downloads/Office/frappe/ecommerce-frappe/logs/redis-queue.error.log
user=yabesh
directory=/home/yabesh/Downloads/Office/frappe/ecommerce-frappe/sites



[program:ecommerce-frappe-redis-socketio]
command=/usr/bin/redis-server /home/yabesh/Downloads/Office/frappe/ecommerce-frappe/config/redis_socketio.conf
priority=1
autostart=true
autorestart=true
stdout_logfile=/home/yabesh/Downloads/Office/frappe/ecommerce-frappe/logs/redis-socketio.log
stderr_logfile=/home/yabesh/Downloads/Office/frappe/ecommerce-frappe/logs/redis-socketio.error.log
user=yabesh
directory=/home/yabesh/Downloads/Office/frappe/ecommerce-frappe/sites



[program:ecommerce-frappe-node-socketio]
command=/home/yabesh/.nvm/versions/node/v18.17.0/bin/node /home/yabesh/Downloads/Office/frappe/ecommerce-frappe/apps/frappe/socketio.js
priority=4
autostart=true
autorestart=true
stdout_logfile=/home/yabesh/Downloads/Office/frappe/ecommerce-frappe/logs/node-socketio.log
stderr_logfile=/home/yabesh/Downloads/Office/frappe/ecommerce-frappe/logs/node-socketio.error.log
user=yabesh
directory=/home/yabesh/Downloads/Office/frappe/ecommerce-frappe


[group:ecommerce-frappe-web]
programs=ecommerce-frappe-frappe-web,ecommerce-frappe-node-socketio



[group:ecommerce-frappe-workers]
programs=ecommerce-frappe-frappe-schedule,ecommerce-frappe-frappe-default-worker,ecommerce-frappe-frappe-short-worker,ecommerce-frappe-frappe-long-worker




[group:ecommerce-frappe-redis]
programs=ecommerce-frappe-redis-cache,ecommerce-frappe-redis-queue,ecommerce-frappe-redis-socketio
